# misc/config/supervisor/position_strategy.ini
[program:position_strategy]
command=bash -lc 'bundle exec unicorn_rails -c config/unicorn.rb -E production -p3000'
numprocs=1
autostart=true
autorestart=true
stopsignal=QUIT
startsecs=1
user=homepage
redirect_stderr=true
stdout_logfile=/var/log/position_strategy/supervisor/stdout.log
stdout_logfile_maxbytes=0
stdout_logfile_backups=0
stderr_logfile=/var/log/position_strategy/supervisor/error.log
stderr_logfile_maxbytes=0
stderr_logfile_backups=0
directory=/home/user-name/position_strategy/current
